# config/systemd/user/battery-warning.timer
[Unit]
Description=Run battery warning periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
