# config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate hyprpaper wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-rotate
